# packaging/upgrade-1.2.0-to-1.3.0/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/mnt/data/logs/pulldb/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; =============================================================================
; MySQL — starts first, highest priority
; Run mysqld directly (not mysqld_safe) so supervisord owns the process
; and can detect crashes and enforce autorestart reliably.
; =============================================================================
[program:mysqld]
command=/usr/sbin/mysqld --user=mysql
autostart=true
autorestart=true
priority=1
startretries=3
startsecs=5
stdout_logfile=/mnt/data/logs/pulldb/mysqld.log
stderr_logfile=/mnt/data/logs/pulldb/mysqld.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3

; =============================================================================
; pulldb-api — waits for MySQL, then starts REST API on :8080
; =============================================================================
[program:pulldb-api]
command=/bin/bash -c "set -a; source /opt/pulldb.service/.env 2>/dev/null; set +a; exec /docker/wait-for-mysql.sh /opt/pulldb.service/venv/bin/pulldb-api"
autostart=true
autorestart=true
priority=10
startretries=10
startsecs=5
user=pulldb_service
environment=HOME="/opt/pulldb.service"
stdout_logfile=/mnt/data/logs/pulldb/pulldb-api.log
stderr_logfile=/mnt/data/logs/pulldb/pulldb-api.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5

; =============================================================================
; pulldb-worker — waits for MySQL, processes restore jobs
; =============================================================================
[program:pulldb-worker]
command=/bin/bash -c "set -a; source /opt/pulldb.service/.env 2>/dev/null; set +a; exec /docker/wait-for-mysql.sh /opt/pulldb.service/venv/bin/pulldb-worker"
autostart=true
autorestart=true
priority=20
startretries=10
startsecs=5
user=pulldb_service
environment=HOME="/opt/pulldb.service"
stdout_logfile=/mnt/data/logs/pulldb/pulldb-worker.log
stderr_logfile=/mnt/data/logs/pulldb/pulldb-worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5

; =============================================================================
; pulldb-web — waits for MySQL, serves Web UI on :8000
; =============================================================================
[program:pulldb-web]
command=/bin/bash -c "set -a; source /opt/pulldb.service/.env 2>/dev/null; set +a; exec /docker/wait-for-mysql.sh /opt/pulldb.service/venv/bin/pulldb-web"
autostart=true
autorestart=true
priority=30
startretries=10
startsecs=5
user=pulldb_service
environment=HOME="/opt/pulldb.service"
stdout_logfile=/mnt/data/logs/pulldb/pulldb-web.log
stderr_logfile=/mnt/data/logs/pulldb/pulldb-web.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
