# packaging/debian/postinst
#!/bin/sh
set -e

INSTALL_PREFIX="/opt/pulldb.service"
SYSTEM_USER="pulldb_service"
SYSTEM_GROUP="pulldb_service"
UNIT_SOURCE="${INSTALL_PREFIX}/scripts/pulldb-worker.service"
UNIT_DEST="/etc/systemd/system/pulldb-worker.service"
ENV_FILE="${INSTALL_PREFIX}/.env"

create_group_and_user() {
    if ! getent group "${SYSTEM_GROUP}" >/dev/null 2>&1; then
        addgroup --system "${SYSTEM_GROUP}"
    fi
    if ! id -u "${SYSTEM_USER}" >/dev/null 2>&1; then
        adduser --system --ingroup "${SYSTEM_GROUP}" --home "${INSTALL_PREFIX}" --shell /usr/sbin/nologin "${SYSTEM_USER}"
    fi
}

ensure_directories() {
    install -d -o "${SYSTEM_USER}" -g "${SYSTEM_GROUP}" -m 0750 "${INSTALL_PREFIX}"
    for dir in logs work scripts dist; do
        install -d -o "${SYSTEM_USER}" -g "${SYSTEM_GROUP}" -m 0750 "${INSTALL_PREFIX}/${dir}"
    done
}

configure_aws_setup() {
    echo "Configuring AWS setup for ${SYSTEM_USER}..."
    USER_HOME=$(getent passwd "${SYSTEM_USER}" | cut -d: -f6)
    
    if [ -z "$USER_HOME" ]; then
        echo "Warning: Could not find home directory for ${SYSTEM_USER}"
        return
    fi

    AWS_DIR="${USER_HOME}/.aws"
    mkdir -p "${AWS_DIR}"
    chown "${SYSTEM_USER}:${SYSTEM_GROUP}" "${AWS_DIR}"
    chmod 700 "${AWS_DIR}"

    # Copy AWS-SETUP.md
    if [ -f "${INSTALL_PREFIX}/AWS-SETUP.md" ]; then
        cp "${INSTALL_PREFIX}/AWS-SETUP.md" "${USER_HOME}/AWS-SETUP.md"
        chown "${SYSTEM_USER}:${SYSTEM_GROUP}" "${USER_HOME}/AWS-SETUP.md"
    fi
    
    echo "AWS setup info copied to ${USER_HOME}/AWS-SETUP.md"
}

validate_install() {
    echo "Validating installation..."
    if [ -x "${INSTALL_PREFIX}/venv/bin/python" ]; then
        if su -s /bin/bash "${SYSTEM_USER}" -c "${INSTALL_PREFIX}/venv/bin/python -c \"import pulldb; print('pullDB module loaded successfully')\""; then
             echo "Validation successful."
        else
             echo "Warning: pullDB module import failed. Check logs and configuration."
        fi
    fi
}

install_systemd_unit() {
    if ! command -v systemctl >/dev/null 2>&1; then
        return
    fi
    if [ -f "${UNIT_SOURCE}" ]; then
        cp "${UNIT_SOURCE}" "${UNIT_DEST}"
        if ! grep -q '^EnvironmentFile=' "${UNIT_DEST}"; then
            sed -i "/^\[Service\]/a EnvironmentFile=${ENV_FILE}" "${UNIT_DEST}"
        fi
        systemctl daemon-reload || true
        systemctl enable pulldb-worker.service || true
    fi
}

case "$1" in
    configure)
        create_group_and_user
        ensure_directories
        
        # Run internal installer to setup venv and install package
        echo "Running internal installer..."
        cd "${INSTALL_PREFIX}"
        if [ -f "scripts/install_pulldb.sh" ]; then
            ./scripts/install_pulldb.sh --prefix "${INSTALL_PREFIX}" --yes --no-systemd --aws-profile dev --secret /pulldb/mysql/coordination-db
        else
            echo "Warning: scripts/install_pulldb.sh not found. Skipping venv setup."
        fi

        # Ensure ownership is correct after installer runs as root
        chown -R "${SYSTEM_USER}:${SYSTEM_GROUP}" "${INSTALL_PREFIX}" || true
        
        install_systemd_unit
        configure_aws_setup
        validate_install
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        ;;
esac

exit 0
