# packaging/debian/postinst
#!/bin/sh
set -e

# =============================================================================
# pullDB Server Package - Post-Installation Script
# =============================================================================
#
# Installation Order:
#   1. Create system user/group
#   2. Merge/create .env and .aws/config (preserving existing values)
#   3. Load variables from .env
#   4. Offer interactive configuration (if terminal available)
#   5. Ensure all paths exist with correct ownership/permissions
#   6. Install Python package and systemd units
#   7. Validate installation
#
# =============================================================================

INSTALL_PREFIX="/opt/pulldb.service"
SYSTEM_USER="pulldb_service"
SYSTEM_GROUP="pulldb_service"
ENV_FILE="${INSTALL_PREFIX}/.env"
MERGE_SCRIPT="${INSTALL_PREFIX}/scripts/merge-config.sh"
CONFIGURE_SCRIPT="${INSTALL_PREFIX}/scripts/configure-pulldb.sh"
BACKUP_TIMESTAMP=$(date +%Y%m%d-%H%M%S)

# Systemd paths
WORKER_UNIT_SOURCE="${INSTALL_PREFIX}/systemd/pulldb-worker.service"
WORKER_UNIT_DEST="/etc/systemd/system/pulldb-worker.service"
API_UNIT_SOURCE="${INSTALL_PREFIX}/systemd/pulldb-api.service"
API_UNIT_DEST="/etc/systemd/system/pulldb-api.service"

# =============================================================================
# Step 1: Create System User and Group
# =============================================================================

create_user_and_group() {
    echo "Step 1: Creating system user and group..."
    
    if ! getent group "${SYSTEM_GROUP}" >/dev/null 2>&1; then
        addgroup --system "${SYSTEM_GROUP}"
        echo "  Created group: ${SYSTEM_GROUP}"
    fi
    
    if ! id -u "${SYSTEM_USER}" >/dev/null 2>&1; then
        adduser --system --ingroup "${SYSTEM_GROUP}" --home "${INSTALL_PREFIX}" --shell /usr/sbin/nologin "${SYSTEM_USER}"
        echo "  Created user: ${SYSTEM_USER}"
    fi
    
    # Ensure base install directory exists
    install -d -o "${SYSTEM_USER}" -g "${SYSTEM_GROUP}" -m 0750 "${INSTALL_PREFIX}"
    for dir in scripts systemd dist logs work after_sql; do
        install -d -o "${SYSTEM_USER}" -g "${SYSTEM_GROUP}" -m 0750 "${INSTALL_PREFIX}/${dir}"
    done
}

# =============================================================================
# Step 2: Merge/Create Configuration Files
# =============================================================================

configure_env_file() {
    echo "Step 2a: Configuring .env file..."
    
    local ENV_EXAMPLE="${INSTALL_PREFIX}/env.example"
    
    if [ -f "${ENV_FILE}" ] && [ -f "${ENV_EXAMPLE}" ]; then
        # Backup existing .env
        cp "${ENV_FILE}" "${ENV_FILE}.backup.${BACKUP_TIMESTAMP}"
        echo "  Backed up existing .env to .env.backup.${BACKUP_TIMESTAMP}"
        
        # Merge existing with new template
        if [ -x "${MERGE_SCRIPT}" ]; then
            "${MERGE_SCRIPT}" env "${ENV_FILE}" "${ENV_EXAMPLE}" "${ENV_FILE}.merged"
            mv "${ENV_FILE}.merged" "${ENV_FILE}"
            echo "  Merged configuration with new template"
        else
            echo "  Warning: Merge script not found, preserving existing .env"
        fi
    elif [ -f "${ENV_EXAMPLE}" ]; then
        # Fresh install - copy template
        cp "${ENV_EXAMPLE}" "${ENV_FILE}"
        echo "  Created .env from template"
    fi
    
    # Secure permissions
    chown "${SYSTEM_USER}:${SYSTEM_GROUP}" "${ENV_FILE}"
    chmod 600 "${ENV_FILE}"
}

configure_aws() {
    echo "Step 2b: Configuring AWS..."
    
    USER_HOME=$(getent passwd "${SYSTEM_USER}" | cut -d: -f6)
    if [ -z "$USER_HOME" ]; then
        echo "  Warning: Could not find home directory for ${SYSTEM_USER}"
        return
    fi

    AWS_DIR="${USER_HOME}/.aws"
    mkdir -p "${AWS_DIR}"
    chown "${SYSTEM_USER}:${SYSTEM_GROUP}" "${AWS_DIR}"
    chmod 700 "${AWS_DIR}"

    # Copy example to .aws directory (always update)
    if [ -f "${INSTALL_PREFIX}/aws.config.example" ]; then
        cp "${INSTALL_PREFIX}/aws.config.example" "${AWS_DIR}/config.example"
        chown "${SYSTEM_USER}:${SYSTEM_GROUP}" "${AWS_DIR}/config.example"
    fi

    # Merge or create .aws/config
    if [ -f "${AWS_DIR}/config" ] && [ -f "${AWS_DIR}/config.example" ]; then
        cp "${AWS_DIR}/config" "${AWS_DIR}/config.backup.${BACKUP_TIMESTAMP}"
        echo "  Backed up existing AWS config"
        
        if [ -x "${MERGE_SCRIPT}" ]; then
            "${MERGE_SCRIPT}" ini "${AWS_DIR}/config" "${AWS_DIR}/config.example" "${AWS_DIR}/config.merged"
            mv "${AWS_DIR}/config.merged" "${AWS_DIR}/config"
            chown "${SYSTEM_USER}:${SYSTEM_GROUP}" "${AWS_DIR}/config"
            chmod 600 "${AWS_DIR}/config"
            echo "  Merged AWS configuration"
        fi
    elif [ -f "${AWS_DIR}/config.example" ]; then
        cp "${AWS_DIR}/config.example" "${AWS_DIR}/config"
        chown "${SYSTEM_USER}:${SYSTEM_GROUP}" "${AWS_DIR}/config"
        chmod 600 "${AWS_DIR}/config"
        echo "  Created AWS config from template"
    fi

    # Copy documentation (only if destination is different from source)
    if [ -f "${INSTALL_PREFIX}/AWS-SETUP.md" ] && [ "${INSTALL_PREFIX}" != "${USER_HOME}" ]; then
        cp "${INSTALL_PREFIX}/AWS-SETUP.md" "${USER_HOME}/AWS-SETUP.md"
        chown "${SYSTEM_USER}:${SYSTEM_GROUP}" "${USER_HOME}/AWS-SETUP.md"
    fi
}

# =============================================================================
# Step 3: Load Variables from .env
# =============================================================================

load_env_variables() {
    echo "Step 3: Loading configuration..."
    
    if [ -f "${ENV_FILE}" ]; then
        # Source the env file to get paths
        set -a
        # shellcheck disable=SC1090
        . "${ENV_FILE}" 2>/dev/null || true
        set +a
        echo "  Loaded ${ENV_FILE}"
    fi
    
    # Set defaults for paths if not defined
    PULLDB_WORK_DIR="${PULLDB_WORK_DIR:-/mnt/data/work/pulldb.service}"
    PULLDB_LOG_DIR="${PULLDB_LOG_DIR:-/mnt/data/logs/pulldb.service}"
    PULLDB_TMP_DIR="${PULLDB_TMP_DIR:-/mnt/data/tmp}"
    PULLDB_CUSTOMERS_AFTER_SQL_DIR="${PULLDB_CUSTOMERS_AFTER_SQL_DIR:-${INSTALL_PREFIX}/after_sql/customer}"
    PULLDB_QA_TEMPLATE_AFTER_SQL_DIR="${PULLDB_QA_TEMPLATE_AFTER_SQL_DIR:-${INSTALL_PREFIX}/after_sql/quality}"
}

# =============================================================================
# Step 4: Interactive Configuration (Optional)
# =============================================================================

offer_interactive_config() {
    # Only offer interactive config if:
    # 1. We have a terminal
    # 2. DEBIAN_FRONTEND is not noninteractive
    # 3. The configure script exists
    
    if [ -t 0 ] && [ "${DEBIAN_FRONTEND:-}" != "noninteractive" ] && [ -x "${CONFIGURE_SCRIPT}" ]; then
        echo ""
        echo "Step 4: Configuration"
        echo ""
        read -p "  Would you like to configure pullDB now? [y/N] " response
        case "$response" in
            [yY][eE][sS]|[yY])
                "${CONFIGURE_SCRIPT}"
                # Reload variables after interactive config
                load_env_variables
                ;;
            *)
                echo "  Skipping interactive configuration."
                echo "  Run 'sudo ${CONFIGURE_SCRIPT}' later to configure."
                ;;
        esac
    else
        # Non-interactive mode: use existing configuration as-is
        echo "Step 4: Using existing configuration (non-interactive mode)"
        if [ -f "${ENV_FILE}" ]; then
            echo "  Preserved: ${ENV_FILE}"
        fi
        if [ -f "${INSTALL_PREFIX}/.aws/config" ]; then
            echo "  Preserved: ${INSTALL_PREFIX}/.aws/config"
        fi
    fi
}

# =============================================================================
# Step 5: Ensure Paths Exist with Correct Permissions
# =============================================================================

ensure_paths() {
    echo "Step 5: Ensuring paths and permissions..."
    
    # Create parent /mnt/data if needed
    if [ ! -d "/mnt/data" ]; then
        mkdir -p /mnt/data
        chmod 755 /mnt/data
    fi
    
    # Work directory
    if [ ! -d "${PULLDB_WORK_DIR}" ]; then
        install -d -o "${SYSTEM_USER}" -g "${SYSTEM_GROUP}" -m 0750 "${PULLDB_WORK_DIR}"
        echo "  Created: ${PULLDB_WORK_DIR}"
    else
        chown "${SYSTEM_USER}:${SYSTEM_GROUP}" "${PULLDB_WORK_DIR}" 2>/dev/null || true
    fi
    
    # Log directory
    if [ ! -d "${PULLDB_LOG_DIR}" ]; then
        install -d -o "${SYSTEM_USER}" -g "${SYSTEM_GROUP}" -m 0750 "${PULLDB_LOG_DIR}"
        echo "  Created: ${PULLDB_LOG_DIR}"
    else
        chown "${SYSTEM_USER}:${SYSTEM_GROUP}" "${PULLDB_LOG_DIR}" 2>/dev/null || true
    fi
    
    # Temp directory (shared, sticky bit)
    if [ ! -d "${PULLDB_TMP_DIR}" ]; then
        install -d -m 1777 "${PULLDB_TMP_DIR}"
        echo "  Created: ${PULLDB_TMP_DIR}"
    fi
    
    # After-SQL directories
    install -d -o "${SYSTEM_USER}" -g "${SYSTEM_GROUP}" -m 0750 "${PULLDB_CUSTOMERS_AFTER_SQL_DIR}"
    install -d -o "${SYSTEM_USER}" -g "${SYSTEM_GROUP}" -m 0750 "${PULLDB_QA_TEMPLATE_AFTER_SQL_DIR}"
    
    # Copy template SQL files if target directories are empty
    TEMPLATE_CUSTOMER="${INSTALL_PREFIX}/template_after_sql/customer"
    TEMPLATE_QA="${INSTALL_PREFIX}/template_after_sql/quality"
    
    if [ -d "${TEMPLATE_CUSTOMER}" ]; then
        if [ -z "$(ls -A "${PULLDB_CUSTOMERS_AFTER_SQL_DIR}" 2>/dev/null)" ]; then
            cp -r "${TEMPLATE_CUSTOMER}"/* "${PULLDB_CUSTOMERS_AFTER_SQL_DIR}/" 2>/dev/null || true
            chown -R "${SYSTEM_USER}:${SYSTEM_GROUP}" "${PULLDB_CUSTOMERS_AFTER_SQL_DIR}"
            echo "  Installed customer after-SQL scripts"
        fi
    fi
    
    if [ -d "${TEMPLATE_QA}" ]; then
        if [ -z "$(ls -A "${PULLDB_QA_TEMPLATE_AFTER_SQL_DIR}" 2>/dev/null)" ]; then
            cp -r "${TEMPLATE_QA}"/* "${PULLDB_QA_TEMPLATE_AFTER_SQL_DIR}/" 2>/dev/null || true
            chown -R "${SYSTEM_USER}:${SYSTEM_GROUP}" "${PULLDB_QA_TEMPLATE_AFTER_SQL_DIR}"
            echo "  Installed QA after-SQL scripts"
        fi
    fi
}

# =============================================================================
# Step 6: Install Python Package and Systemd Units
# =============================================================================

install_python_package() {
    echo "Step 6a: Installing Python package..."
    
    VENV_DIR="${INSTALL_PREFIX}/venv"
    
    if [ -d "${VENV_DIR}" ]; then
        echo "  Existing venv found, upgrading..."
        # Clean up any corrupted distributions (pip temp files starting with ~)
        find "${VENV_DIR}/lib" -type d -name '~*' -exec rm -rf {} + 2>/dev/null || true
    else
        echo "  Creating new virtual environment..."
        python3 -m venv "${VENV_DIR}"
    fi
    
    # Install/upgrade the wheel
    if ls "${INSTALL_PREFIX}/dist"/pulldb-*.whl 1>/dev/null 2>&1; then
        "${VENV_DIR}/bin/pip" install --upgrade pip wheel >/dev/null
        "${VENV_DIR}/bin/pip" install --upgrade --force-reinstall "${INSTALL_PREFIX}/dist"/pulldb-*.whl
        echo "  pulldb package installed successfully"
    else
        echo "  Warning: No wheel found in ${INSTALL_PREFIX}/dist/"
    fi
    
    # Ensure ownership
    chown -R "${SYSTEM_USER}:${SYSTEM_GROUP}" "${INSTALL_PREFIX}" || true
}

check_database_schema() {
    echo "Step 6b: Checking database schema..."
    
    # Check if pulldb_service database exists with required tables
    if sudo mysql -N -e "SELECT 1 FROM information_schema.tables WHERE table_schema='pulldb_service' AND table_name='jobs'" 2>/dev/null | grep -q 1; then
        echo "  Database schema present"
    else
        echo "  Warning: Database schema not found"
        echo "  v0.0.6 requires fresh schema setup. Run:"
        echo "    sudo mysql -e 'CREATE DATABASE IF NOT EXISTS pulldb_service'"
        echo "    cat ${INSTALL_PREFIX}/schema/pulldb_service/*.sql | sudo mysql pulldb_service"
    fi
}

install_systemd_units() {
    echo "Step 6c: Installing systemd units..."
    
    if ! command -v systemctl >/dev/null 2>&1; then
        echo "  systemd not available, skipping"
        return
    fi
    
    # Remove old copies if they exist (migrate to symlinks)
    [ -f "${WORKER_UNIT_DEST}" ] && [ ! -L "${WORKER_UNIT_DEST}" ] && rm -f "${WORKER_UNIT_DEST}"
    [ -f "${API_UNIT_DEST}" ] && [ ! -L "${API_UNIT_DEST}" ] && rm -f "${API_UNIT_DEST}"
    
    # Create symlinks
    if [ -f "${WORKER_UNIT_SOURCE}" ]; then
        ln -sf "${WORKER_UNIT_SOURCE}" "${WORKER_UNIT_DEST}"
        systemctl enable pulldb-worker.service >/dev/null 2>&1 || true
        echo "  Enabled pulldb-worker.service"
    fi
    
    if [ -f "${API_UNIT_SOURCE}" ]; then
        ln -sf "${API_UNIT_SOURCE}" "${API_UNIT_DEST}"
        systemctl enable pulldb-api.service >/dev/null 2>&1 || true
        echo "  Enabled pulldb-api.service"
    fi
    
    systemctl daemon-reload || true
}

start_services() {
    echo "Step 6d: Starting services..."
    
    if ! command -v systemctl >/dev/null 2>&1; then
        echo "  systemd not available, skipping"
        return
    fi
    
    # Start API service first (worker depends on it being available)
    if systemctl is-enabled pulldb-api.service >/dev/null 2>&1; then
        systemctl restart pulldb-api.service 2>/dev/null || true
        if systemctl is-active --quiet pulldb-api.service 2>/dev/null; then
            echo "  Started pulldb-api.service"
        else
            echo "  Warning: pulldb-api.service failed to start"
        fi
    fi
    
    # Start worker service
    if systemctl is-enabled pulldb-worker.service >/dev/null 2>&1; then
        systemctl restart pulldb-worker.service 2>/dev/null || true
        if systemctl is-active --quiet pulldb-worker.service 2>/dev/null; then
            echo "  Started pulldb-worker.service"
        else
            echo "  Warning: pulldb-worker.service failed to start"
        fi
    fi
}

# =============================================================================
# Step 7: Validate Installation
# =============================================================================

validate_installation() {
    echo "Step 7: Validating installation..."
    
    if [ -x "${INSTALL_PREFIX}/venv/bin/python" ]; then
        if su -s /bin/bash "${SYSTEM_USER}" -c "${INSTALL_PREFIX}/venv/bin/python -c \"import pulldb; print('  pullDB module loaded successfully')\"" 2>/dev/null; then
            echo "  Validation successful"
        else
            echo "  Warning: pullDB module import failed"
        fi
    fi
}

# =============================================================================
# Completion Message
# =============================================================================

print_completion() {
    echo ""
    echo "════════════════════════════════════════════════════════════════"
    echo "  Installation Complete!"
    echo "════════════════════════════════════════════════════════════════"
    echo ""
    echo "  Configuration: ${ENV_FILE}"
    echo "  Configure:     sudo ${CONFIGURE_SCRIPT}"
    echo ""
    echo "  Manage services:"
    echo "    sudo systemctl restart pulldb-api pulldb-worker"
    echo ""
    echo "  Check status:"
    echo "    sudo systemctl status pulldb-api pulldb-worker"
    echo ""
}

# =============================================================================
# Main Entry Point
# =============================================================================

case "$1" in
    configure)
        create_user_and_group
        configure_env_file
        configure_aws
        load_env_variables
        offer_interactive_config
        ensure_paths
        install_python_package
        check_database_schema
        install_systemd_units
        start_services
        validate_installation
        print_completion
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        ;;
esac

exit 0
