# packaging/debian/postrm
#!/bin/sh
set -e

# =============================================================================
# pullDB Server Package - Post-Removal Script
# =============================================================================
#
# dpkg -r pulldb (remove): Removes everything EXCEPT .env and .aws/config
# dpkg -P pulldb (purge):  Removes EVERYTHING including config files
#
# =============================================================================

INSTALL_PREFIX="/opt/pulldb.service"
SYSTEM_USER="pulldb_service"
WORKER_UNIT_DEST="/etc/systemd/system/pulldb-worker.service"
WORKER_TEMPLATE_DEST="/etc/systemd/system/pulldb-worker@.service"
API_UNIT_DEST="/etc/systemd/system/pulldb-api.service"
WEB_UNIT_DEST="/etc/systemd/system/pulldb-web.service"
RETENTION_SERVICE_DEST="/etc/systemd/system/pulldb-retention.service"
RETENTION_TIMER_DEST="/etc/systemd/system/pulldb-retention.timer"

cleanup_systemd() {
    echo "Removing systemd units..."
    if command -v systemctl >/dev/null 2>&1; then
        systemctl stop pulldb-worker.service >/dev/null 2>&1 || true
        systemctl disable pulldb-worker.service >/dev/null 2>&1 || true
        # Stop and disable any running template instances (worker@1, worker@2, etc.)
        systemctl stop 'pulldb-worker@*.service' >/dev/null 2>&1 || true
        systemctl disable 'pulldb-worker@*.service' >/dev/null 2>&1 || true
        systemctl stop pulldb-api.service >/dev/null 2>&1 || true
        systemctl disable pulldb-api.service >/dev/null 2>&1 || true
        systemctl stop pulldb-web.service >/dev/null 2>&1 || true
        systemctl disable pulldb-web.service >/dev/null 2>&1 || true
        systemctl stop pulldb-retention.timer >/dev/null 2>&1 || true
        systemctl disable pulldb-retention.timer >/dev/null 2>&1 || true
        
        # Remove main unit files
        rm -f "${WORKER_UNIT_DEST}" "${WORKER_TEMPLATE_DEST}" "${API_UNIT_DEST}" "${WEB_UNIT_DEST}" || true
        rm -f "${RETENTION_SERVICE_DEST}" "${RETENTION_TIMER_DEST}" || true
        
        # Remove template instance symlinks (pulldb-worker@1.service, @2, @3, etc.)
        rm -f /etc/systemd/system/pulldb-worker@[0-9]*.service 2>/dev/null || true
        
        # Remove any service override directories
        rm -rf /etc/systemd/system/pulldb-worker.service.d 2>/dev/null || true
        rm -rf /etc/systemd/system/pulldb-api.service.d 2>/dev/null || true
        rm -rf /etc/systemd/system/pulldb-web.service.d 2>/dev/null || true
        rm -rf /etc/systemd/system/pulldb-retention.service.d 2>/dev/null || true
        
        # Remove multi-user.target.wants symlinks
        rm -f /etc/systemd/system/multi-user.target.wants/pulldb-*.service 2>/dev/null || true
        rm -f /etc/systemd/system/multi-user.target.wants/pulldb-*.timer 2>/dev/null || true
        
        systemctl daemon-reload || true
        systemctl reset-failed 2>/dev/null || true
    else
        rm -f "${WORKER_UNIT_DEST}" "${WORKER_TEMPLATE_DEST}" "${API_UNIT_DEST}" "${WEB_UNIT_DEST}" || true
        rm -f "${RETENTION_SERVICE_DEST}" "${RETENTION_TIMER_DEST}" || true
        rm -f /etc/systemd/system/pulldb-worker@[0-9]*.service 2>/dev/null || true
        rm -rf /etc/systemd/system/pulldb-*.service.d 2>/dev/null || true
    fi
}

# Remove application files but preserve .env and .aws/config
remove_keeping_config() {
    echo "Removing pullDB (preserving configuration)..."
    
    # Remove CLI tools from /usr/local/bin
    rm -f /usr/local/bin/pulldb 2>/dev/null || true
    rm -f /usr/local/bin/pulldb-admin 2>/dev/null || true
    echo "  Removed: /usr/local/bin/pulldb, /usr/local/bin/pulldb-admin"
    
    # Remove sudoers rule for pulldb-admin
    rm -f /etc/sudoers.d/pulldb-admin 2>/dev/null || true
    echo "  Removed: /etc/sudoers.d/pulldb-admin"
    
    # Remove pulldb CA certificate from system trust store
    if [ -f /usr/local/share/ca-certificates/pulldb-service.crt ]; then
        rm -f /usr/local/share/ca-certificates/pulldb-service.crt
        update-ca-certificates 2>/dev/null || true
        echo "  Removed: CA certificate from system trust store"
    fi
    
    if [ ! -d "${INSTALL_PREFIX}" ]; then
        return 0
    fi
    
    # Remove all known directories
    rm -rf "${INSTALL_PREFIX}/venv" 2>/dev/null || true
    rm -rf "${INSTALL_PREFIX}/dist" 2>/dev/null || true
    rm -rf "${INSTALL_PREFIX}/scripts" 2>/dev/null || true
    rm -rf "${INSTALL_PREFIX}/systemd" 2>/dev/null || true
    rm -rf "${INSTALL_PREFIX}/logs" 2>/dev/null || true
    rm -rf "${INSTALL_PREFIX}/work" 2>/dev/null || true
    rm -rf "${INSTALL_PREFIX}/after_sql" 2>/dev/null || true
    rm -rf "${INSTALL_PREFIX}/template_after_sql" 2>/dev/null || true
    
    # Remove known files (but NOT .env)
    rm -f "${INSTALL_PREFIX}/AWS-SETUP.md" 2>/dev/null || true
    rm -f "${INSTALL_PREFIX}/SERVICE-README.md" 2>/dev/null || true
    rm -f "${INSTALL_PREFIX}/env.example" 2>/dev/null || true
    rm -f "${INSTALL_PREFIX}/aws.config.example" 2>/dev/null || true
    
    # Clean .aws directory but preserve config
    if [ -d "${INSTALL_PREFIX}/.aws" ]; then
        # Remove everything except config
        find "${INSTALL_PREFIX}/.aws" -type f ! -name 'config' -delete 2>/dev/null || true
        find "${INSTALL_PREFIX}/.aws" -mindepth 1 -type d -empty -delete 2>/dev/null || true
    fi
    
    # Remove .env backups but keep .env
    find "${INSTALL_PREFIX}" -maxdepth 1 -name '.env.backup.*' -delete 2>/dev/null || true
    find "${INSTALL_PREFIX}" -maxdepth 1 -name '.env.merged' -delete 2>/dev/null || true
    
    # Check what's left
    if [ -f "${INSTALL_PREFIX}/.env" ]; then
        echo "  Preserved: ${INSTALL_PREFIX}/.env"
    fi
    if [ -f "${INSTALL_PREFIX}/.aws/config" ]; then
        echo "  Preserved: ${INSTALL_PREFIX}/.aws/config"
    fi
    
    echo "  To remove completely: sudo dpkg -P pulldb"
}

# Remove everything including config and user account
purge_all() {
    echo "Purging pullDB (removing all files, configuration, and user account)..."
    
    # Remove CLI tools from /usr/local/bin
    rm -f /usr/local/bin/pulldb 2>/dev/null || true
    rm -f /usr/local/bin/pulldb-admin 2>/dev/null || true
    
    # Remove sudoers rule for pulldb-admin
    rm -f /etc/sudoers.d/pulldb-admin 2>/dev/null || true
    
    # Remove pulldb CA certificate from system trust store
    if [ -f /usr/local/share/ca-certificates/pulldb-service.crt ]; then
        rm -f /usr/local/share/ca-certificates/pulldb-service.crt
        update-ca-certificates 2>/dev/null || true
        echo "  Removed: CA certificate from system trust store"
    fi
    
    # Remove the entire install directory
    if [ -d "${INSTALL_PREFIX}" ]; then
        rm -rf "${INSTALL_PREFIX}"
        echo "  Removed: ${INSTALL_PREFIX}"
    fi
    
    # Remove system user and group
    if id -u "${SYSTEM_USER}" >/dev/null 2>&1; then
        deluser --system "${SYSTEM_USER}" 2>/dev/null || true
        echo "  Removed user: ${SYSTEM_USER}"
    fi
    
    if getent group "${SYSTEM_USER}" >/dev/null 2>&1; then
        delgroup --system "${SYSTEM_USER}" 2>/dev/null || true
        echo "  Removed group: ${SYSTEM_USER}"
    fi
    
    echo "  Purge complete"
}

case "$1" in
    remove)
        cleanup_systemd
        remove_keeping_config
        ;;
    purge)
        cleanup_systemd
        purge_all
        ;;
    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        ;;
    *)
        ;;
esac

exit 0
