# packaging/debian/postinst
#!/bin/sh
set -e

INSTALL_PREFIX="/opt/pulldb"
SYSTEM_USER="pulldb"
SYSTEM_GROUP="pulldb"
UNIT_SOURCE="${INSTALL_PREFIX}/scripts/pulldb-worker.service"
UNIT_DEST="/etc/systemd/system/pulldb-worker.service"
ENV_FILE="${INSTALL_PREFIX}/.env"

create_group_and_user() {
    if ! getent group "${SYSTEM_GROUP}" >/dev/null 2>&1; then
        addgroup --system "${SYSTEM_GROUP}"
    fi
    if ! id -u "${SYSTEM_USER}" >/dev/null 2>&1; then
        adduser --system --ingroup "${SYSTEM_GROUP}" --home "${INSTALL_PREFIX}" --shell /usr/sbin/nologin "${SYSTEM_USER}"
    fi
}

ensure_directories() {
    install -d -o "${SYSTEM_USER}" -g "${SYSTEM_GROUP}" -m 0750 "${INSTALL_PREFIX}"
    for dir in logs work scripts; do
        install -d -o "${SYSTEM_USER}" -g "${SYSTEM_GROUP}" -m 0750 "${INSTALL_PREFIX}/${dir}"
    done
}

write_default_env() {
    if [ -f "${ENV_FILE}" ]; then
        return
    fi
    cat > "${ENV_FILE}" <<EOF
# pullDB runtime configuration
# Generated by packaging maintainer scripts
PULLDB_INSTALL_PREFIX=${INSTALL_PREFIX}
PULLDB_AWS_PROFILE=dev
PULLDB_COORDINATION_SECRET=/pulldb/mysql/coordination-db
EOF
    chown "${SYSTEM_USER}:${SYSTEM_GROUP}" "${ENV_FILE}"
    chmod 0640 "${ENV_FILE}"
}

install_systemd_unit() {
    if ! command -v systemctl >/dev/null 2>&1; then
        return
    fi
    if [ -f "${UNIT_SOURCE}" ]; then
        cp "${UNIT_SOURCE}" "${UNIT_DEST}"
        if ! grep -q '^EnvironmentFile=' "${UNIT_DEST}"; then
            sed -i "/^\[Service\]/a EnvironmentFile=${ENV_FILE}" "${UNIT_DEST}"
        fi
        systemctl daemon-reload || true
        systemctl enable pulldb-worker.service || true
    fi
}

case "$1" in
    configure)
        create_group_and_user
        ensure_directories
        chown -R "${SYSTEM_USER}:${SYSTEM_GROUP}" "${INSTALL_PREFIX}" || true
        write_default_env
        install_systemd_unit
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        ;;
esac

exit 0
