# packaging/debian/prerm
#!/bin/sh
set -e

INSTALL_PREFIX="/opt/pulldb.service"
ENV_FILE="${INSTALL_PREFIX}/.env"
# Safe backup location that survives purge (outside INSTALL_PREFIX)
ENV_BACKUP_DIR="/var/lib/pulldb"
ENV_BACKUP_FILE="${ENV_BACKUP_DIR}/.env.pre-purge"

WORKER_SERVICE="pulldb-worker.service"
API_SERVICE="pulldb-api.service"
WEB_SERVICE="pulldb-web.service"
RETENTION_TIMER="pulldb-retention.timer"

stop_services() {
    if command -v systemctl >/dev/null 2>&1; then
        # Stop retention timer first
        systemctl stop "${RETENTION_TIMER}" >/dev/null 2>&1 || true
        systemctl disable "${RETENTION_TIMER}" >/dev/null 2>&1 || true
        # Stop any worker template instances (worker@1, worker@2, etc.)
        systemctl stop 'pulldb-worker@*.service' >/dev/null 2>&1 || true
        systemctl stop "${WORKER_SERVICE}" >/dev/null 2>&1 || true
        systemctl disable "${WORKER_SERVICE}" >/dev/null 2>&1 || true
        systemctl stop "${API_SERVICE}" >/dev/null 2>&1 || true
        systemctl disable "${API_SERVICE}" >/dev/null 2>&1 || true
        systemctl stop "${WEB_SERVICE}" >/dev/null 2>&1 || true
        systemctl disable "${WEB_SERVICE}" >/dev/null 2>&1 || true
    fi
}

backup_env_for_purge() {
    # Backup .env to a location OUTSIDE INSTALL_PREFIX before purge
    # This allows postinst to restore settings even after a full purge
    if [ -f "${ENV_FILE}" ]; then
        mkdir -p "${ENV_BACKUP_DIR}"
        cp "${ENV_FILE}" "${ENV_BACKUP_FILE}"
        chmod 600 "${ENV_BACKUP_FILE}"
        echo "  Backed up .env to ${ENV_BACKUP_FILE} (will be restored on reinstall)"
    fi
}

case "$1" in
    remove|upgrade|deconfigure)
        stop_services
        # For remove/upgrade, also backup .env in case user later purges
        backup_env_for_purge
        ;;
    failed-upgrade)
        ;;
    *)
        ;;
esac

exit 0
