# packaging/debian/prerm
#!/bin/sh
set -e

WORKER_SERVICE="pulldb-worker.service"
API_SERVICE="pulldb-api.service"
WEB_SERVICE="pulldb-web.service"

stop_services() {
    if command -v systemctl >/dev/null 2>&1; then
        systemctl stop "${WORKER_SERVICE}" >/dev/null 2>&1 || true
        systemctl disable "${WORKER_SERVICE}" >/dev/null 2>&1 || true
        systemctl stop "${API_SERVICE}" >/dev/null 2>&1 || true
        systemctl disable "${API_SERVICE}" >/dev/null 2>&1 || true
        systemctl stop "${WEB_SERVICE}" >/dev/null 2>&1 || true
        systemctl disable "${WEB_SERVICE}" >/dev/null 2>&1 || true
    fi
}

case "$1" in
    remove|upgrade|deconfigure)
        stop_services
        ;;
    failed-upgrade)
        ;;
    *)
        ;;
esac

exit 0
