# packaging/systemd/pulldb-web.service
[Unit]
Description=pullDB Web UI Service (Web dashboard on port 8000)
After=network.target mysql.service pulldb-api.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/pulldb.service
EnvironmentFile=/opt/pulldb.service/.env
# Web-only mode: use dedicated entry point on port 8000
Environment=PULLDB_WEB_PORT=8000
ExecStart=/opt/pulldb.service/venv/bin/pulldb-web
Restart=on-failure
RestartSec=5
User=pulldb_service
Group=pulldb_service
KillSignal=SIGTERM
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
