# scripts/pulldb-worker.service
; Example systemd unit file for pullDB worker daemon
; Install path: /etc/systemd/system/pulldb-worker.service
; Enable: systemctl enable pulldb-worker.service
; Start:  systemctl start pulldb-worker.service
; Logs captured via structured JSON (consider journalctl -u pulldb-worker -o cat)

[Unit]
Description=pullDB Restore Worker Daemon
After=network.target mysql.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
# Secrets + MySQL live in dev account: use pr-dev profile for credential resolution
Environment=PULLDB_AWS_PROFILE=pr-dev
# AWS_PROFILE remains unset so the EC2 instance profile handles S3 access.
WorkingDirectory=/opt/pulldb.service
ExecStart=/usr/bin/env pulldb-worker
Restart=on-failure
RestartSec=5
User=pulldb
Group=pulldb
KillSignal=SIGTERM
TimeoutStopSec=30
# Hardening (extend as deployment matures)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
